# cookbooks/flask_app/files/default/demo_flask_app.conf
# demo_flask_app generated by chef

description "demo_flask_app"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

script
  cd /opt/demo_flask_app
  /usr/bin/python27 -m gunicorn.app.wsgiapp --bind unix:/opt/demo_flask_app/app.sock main
end script
